# rtorrent.service
[Unit]
Description=rTorrent (Otimizado Daemon)
After=network.target mnt-sdb2.mount
Requires=mnt-sdb2.mount
ConditionPathIsMountPoint=/mnt/sdb2

[Service]
Type=simple
User=pi
Group=pi
UMask=0002
WorkingDirectory=/home/pi
ExecStart=/usr/bin/rtorrent -n -o import=/home/pi/.rtorrent.rc
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
KillMode=process
Nice=-5

[Install]
WantedBy=multi-user.target
